# init/ssh-container-on-demand.socket
[Socket]
ListenStream=0.0.0.0:2222

[Install]
WantedBy=sockets.target
